# services/fuseki-e4s-v3.service
[Unit]
Description=Apache Fuseki Server - E4s Standard v3
After=network.target

[Service]
Type=simple
User=fuseki
Environment=FUSEKI_BASE=/data/fuseki/databases
Environment=FUSEKI_HOME=/opt/fuseki
Environment=JAVA_OPTS=-Xms8g -Xmx32g
ExecStart=/opt/fuseki/fuseki-server --config=/opt/fuseki/configuration/config.ttl --port=3030
Restart=on-failure

[Install]
WantedBy=multi-user.target
